# docker/pkg/sota-client.conf
#/etc/init/sota-client.conf

description "Start SOTA Client on startup."

start on startup
stop on shutdown

exec /opt/ats/sota_client/sota_client -c /opt/ats/sota_client/client.toml
